# .config/systemd/user/wallpaper.service
[Unit]
Description=Change wallpaper

[Service]
ExecStart=%h/.local/bin/changepaper.sh
Type=oneshot

[Install]
WantedBy=default.target
